# scripts/podman/shehzadalgo.container.in
# ShehzadAlgo gateway — Podman Quadlet (rootless)
# Installed by setup-podman.sh into shehzadalgo's ~/.config/containers/systemd/
# {{SHEHZADALGO_HOME}} is replaced at install time.

[Unit]
Description=ShehzadAlgo gateway (rootless Podman)

[Container]
Image=shehzadalgo:local
ContainerName=shehzadalgo
UserNS=keep-id
Volume={{SHEHZADALGO_HOME}}/.shehzadalgo:/home/node/.shehzadalgo
EnvironmentFile={{SHEHZADALGO_HOME}}/.shehzadalgo/.env
Environment=HOME=/home/node
Environment=TERM=xterm-256color
PublishPort=18789:18789
PublishPort=18790:18790
Pull=never
Exec=node dist/index.js gateway --bind lan --port 18789

[Service]
TimeoutStartSec=300
Restart=on-failure

[Install]
WantedBy=default.target
